# ipsec_server/system/certbot_renew.timer
[Unit]
Description=Weekly renewal of lets encrypt certificate

[Timer]
OnCalendar=Wed 06:30:00
Persistent=true

[Install]
WantedBy=timers.target
